# systemd/system/aur-repo.target
[Unit]
Description=Target for all AUR packages
StopWhenUnneeded=yes

[Install]
WantedBy=multi-user.target
